# .config/autostart/ChromeKiosk.desktop
[Desktop Entry] 
Type=Application
Exec=/home/pi/start_chromium.sh

