# .config/autostart/sound.desktop
[Desktop Entry] 
Type=Application
Exec=/home/pi/start_sound.sh

